# three-tier-app/ansible/Roles/application/templates/app1.service.j2
[Unit]
Description=Application One

[Service]
ExecStart=/usr/bin/python3 /vagrant/app_one/app_one.py
Restart=always

[Install]
WantedBy=default.target