# three-tier-app/ansible/Roles/application/templates/app2.service.j2
[Unit]
Description=Application Two

[Service]
ExecStart=/usr/bin/python3 /vagrant/app_two/app_two.py
Restart=always

[Install]
WantedBy=default.target